# optimus-api.service
[Unit]
Description=optimus-api
After=network.target

[Service]
User=opc
WorkingDirectory=/home/opc/projects/optimus
EnvironmentFile=-/home/opc/projects/optimus/.env
ExecStart=/home/opc/projects/optimus/.venv/bin/gunicorn app.main:app \
    -w 1 -k uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
